# files/etc/uci-defaults/99-init-settings
uci set luci.main.mediaurlbase='/luci-static/argon'

exit 0
